# platforms/linux/initfiles/systemd/git-auto-deploy.service
[Unit]
Description=GitAutoDeploy

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/Git-Auto-Deploy/
ExecStart=/usr/bin/python /opt/Git-Auto-Deploy/gitautodeploy --daemon-mode

[Install]
WantedBy=multi-user.target
